# deploy/race-technik/launchagents/com.racetechnik.nightly-github-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.racetechnik.nightly-github-sync</string>

    <!--
        Daily at 23:00 SAST (local time).
        Mac Mini timezone must be Africa/Johannesburg.
    -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/raceai/.amalfiai/workspace/scripts/nightly-github-sync.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/raceai</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/raceai/.amalfiai/workspace/out/nightly-github-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/raceai/.amalfiai/workspace/out/nightly-github-sync.err.log</string>

</dict>
</plist>
