# deploy/race-technik/launchagents/com.racetechnik.sophia-cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.racetechnik.sophia-cron</string>

    <!-- Poll every 1200 seconds (20 minutes) -->
    <key>StartInterval</key>
    <integer>1200</integer>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/raceai/.amalfiai/workspace/scripts/run-claude-job.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/raceai</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/raceai/.amalfiai/workspace/out/sophia-cron.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/raceai/.amalfiai/workspace/out/sophia-cron.err.log</string>

</dict>
</plist>
