# launchagents/ai.openclaw.token-keepalive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.token-keepalive</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/henryburton/.openclaw/scripts/claude-token-keepalive.sh</string>
  </array>

  <!-- Every 4 hours. Tokens last ~7h so this guarantees at least one check before expiry. -->
  <key>StartInterval</key>
  <integer>14400</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/henryburton/.openclaw/logs/token-keepalive.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/henryburton/.openclaw/logs/token-keepalive.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/henryburton</string>
  </dict>
</dict>
</plist>
