# launchagents/com.amalfiai.weekly-client-reports.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.amalfiai.weekly-client-reports</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/henryburton/.openclaw/workspace-anthropic/scripts/weekly-reports/generate-weekly-client-reports.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>2</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/henryburton/.openclaw/workspace-anthropic/logs/weekly-client-reports.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/henryburton/.openclaw/workspace-anthropic/logs/weekly-client-reports.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>/Users/henryburton</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
