# launchagents/com.amalfiai.finance-poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.amalfiai.finance-poller</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/henryburton/.openclaw/workspace-anthropic/scripts/finance-poller.mjs</string>
  </array>

  <!-- Run at 07:00 and 18:00 SAST (local system time) -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>7</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/henryburton/.openclaw/workspace-anthropic/out/finance-poller.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/henryburton/.openclaw/workspace-anthropic/out/finance-poller.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/henryburton</string>
  </dict>
</dict>
</plist>
