# launchagents/com.amalfiai.salah-morning-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.amalfiai.salah-morning-brief</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/henryburton/.openclaw/workspace-anthropic/scripts/salah-morning-brief.sh</string>
    </array>

    <!-- 08:00 SAST — change Hour to Salah's preferred time once confirmed -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/henryburton/.openclaw/workspace-anthropic/out/salah-morning-brief.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/henryburton/.openclaw/workspace-anthropic/out/salah-morning-brief.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>HOME</key>
      <string>/Users/henryburton</string>
    </dict>
  </dict>
</plist>
