# clients/vanta-studios/launchagents/com.amalfiai.vanta-lead-verify.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.amalfiai.vanta-lead-verify</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/henryburton/.openclaw/workspace-anthropic/scripts/vanta-lead-verify.sh</string>
    </array>

    <!-- 10:00 SAST daily (after discovery at 09:00) -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>10</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/henryburton/.openclaw/workspace-anthropic/out/vanta-lead-verify.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/henryburton/.openclaw/workspace-anthropic/out/vanta-lead-verify.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>HOME</key>
      <string>/Users/henryburton</string>
    </dict>
  </dict>
</plist>
